# root/etc/systemd/system/mpvsync.service
[Unit]
Description=Unused old MPV sync service

[Service]
Type=oneshot
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
